# .config/systemd/user/plank@.service
[Unit]
Description=Plank Dock %i

[Service]
ExecStart=/usr/bin/plank -n dock%i
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
