# .config/systemd/user/mpd.service.d/override.conf
[Service]
# Delay units that depend on this, to make sure mpd state is set
# This helps to avoid playing/starting notification on startup
ExecStartPost=/usr/bin/sleep 3
