# .config/autostart/syncthing.desktop
[Desktop Entry]
Name=Syncthing
Exec=syncthing --no-browser
Terminal=false
Type=Application
