# .config/systemd/user/mpd-notification.service.d/override.conf
[Unit]
After=xfce4-notifyd.service

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
